# scripts/local.porter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Example LaunchAgent for a *source checkout* only. End users should use Porter.app Set & forget. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.porter.dev</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/ABSOLUTE/PATH/TO/porter/scripts/start-porter.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/porter.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/porter.launchd.err.log</string>
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/porter</string>
</dict>
</plist>
